# scripts/com.claude.timer-keeper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.timer-keeper</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/juanvergara/Desktop/Otros/otros/evasystem/scripts/claude-timer-keeper.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/claude-timer-keeper.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/claude-timer-keeper.err</string>
</dict>
</plist>
